# systemd/overrides/llm-stack.service.d/docker-wait.conf
[Service]
# Wait up to 30s for docker daemon before attempting compose up
ExecStartPre=
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do docker info >/dev/null 2>&1 && exit 0; sleep 1; done; echo "Docker not ready after 30s"; exit 1'
ExecStartPre=/home/hapax/llm-stack/generate-env.sh
