# systemd/units/hapax-dmn.service.d/gpu-pin.conf
[Service]
# LRR Phase 3 item 1 — migrate hapax-dmn from GPU 1 (3090) to GPU 0
# (5060 Ti) to free the 3090 for Hermes 3 70B.
#
# Under Option α (pre-Phase 3), hapax-dmn runs the daimonion voice
# stack on GPU 1 (~3.4 GiB: faster-whisper STT + embedding). GPU 1
# also hosts the current TabbyAPI Qwen3.5-9B (~5.7 GiB). Combined
# usage is ~9.4 GiB of 23.6 GiB — comfortable.
#
# Under Option γ, the 3090 must host Hermes 3 70B's primary weight
# slice (~23.5 GiB + activations + KV cache). There is NO room for
# hapax-dmn on GPU 1 once Hermes 3 is loaded. The daimonion moves
# to GPU 0 (5060 Ti, 15.5 GiB total) alongside the compositor (3.3)
# and hapax-imagination (0.3). New GPU 0 budget:
#
#   compositor      3.3 GiB
#   imagination     0.3
#   hapax-dmn       3.4   ← moved here
#   Hermes overflow 2.75  ← Option γ slice
#   faster-whisper  2.5   ← already on GPU 0 via compositor venv
#   ─────────────────────
#   total          12.25 GiB   (3.25 GiB headroom on 16 GiB card)
#
# Beta verified the math in the supplement drop (see
# ~/.cache/hapax/relay/context/2026-04-14-beta-phase-3-supplement-
# verified-preconditions.md § 1.4). The 3.25 GiB headroom leaves
# room for Bundle 6's optional Kokoro-GPU migration if the operator
# chooses to pursue it.
#
# `CUDA_DEVICE_ORDER=PCI_BUS_ID` is REQUIRED so `CUDA_VISIBLE_DEVICES=0`
# resolves to the 5060 Ti and not to whichever card CUDA's default
# FASTEST_FIRST heuristic decides is "fastest". See the companion
# tabbyapi.service.d/gpu-pin.conf drop-in for the same ordering note.
#
# Rollback to Option α: replace `CUDA_VISIBLE_DEVICES=0` with
# `CUDA_VISIBLE_DEVICES=1` and restart.
Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
Environment="CUDA_VISIBLE_DEVICES=0"
