# systemd/units/contact-mic-recorder.service.d/archive-path.conf
[Service]
# LRR Phase 2 item 5 — Cortado FLAC output path migration.
# See the audio-recorder drop-in for rationale.
Environment=HAPAX_AUDIO_ARCHIVE_ROOT=%h/hapax-state/stream-archive/audio
