# systemd/overrides/studio-fx-output.service.d/priority.conf
[Service]
# Tier 2: Background — low-rate snapshot feed, not latency-sensitive
Nice=10
CPUWeight=25
IOWeight=50
