# systemd/overrides/dev/daily-briefing.timer
[Timer]
OnCalendar=
OnCalendar=*-*-* 03/4:00:00
