# systemd/overrides/dev/digest.timer
[Timer]
OnCalendar=
OnCalendar=*-*-* 00/2:45:00
