# systemd/overrides/hapax-stack.service.d/ordering.conf
[Unit]
After=llm-stack.service logos-api.service pipewire.service
Wants=llm-stack.service
