# systemd/overrides/hapax-stack.service.d/priority.conf
[Service]
# Tier 1: Interactive voice/perception — high priority but yields to compositor
Nice=0
CPUWeight=200
IOWeight=200
