# systemd/overrides/rnd/knowledge-maint.timer
[Timer]
OnCalendar=
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec=5min
