# systemd/overrides/rnd/profile-update.timer
[Timer]
OnBootSec=
OnUnitActiveSec=
OnBootSec=5min
OnUnitActiveSec=45min
RandomizedDelaySec=5min
