# systemd/units-pi6/chrome-sync.service
[Unit]
Description=Chrome history + bookmarks sync (Pi6 offload)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/hapax/bin/sync-agent-wrapper.sh chrome_sync .config/google-chrome/Default/History .config/google-chrome/Default/Bookmarks
MemoryMax=512M
CPUQuota=50%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chrome-sync

[Install]
WantedBy=default.target
