# systemd/units-pi6/claude-code-sync.service
[Unit]
Description=Claude Code transcript sync (Pi6 offload)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/hapax/bin/sync-agent-wrapper.sh claude_code_sync .claude/projects/
MemoryMax=512M
CPUQuota=50%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-code-sync

[Install]
WantedBy=default.target
