# systemd/units-pi6/langfuse-sync.service
[Unit]
Description=Langfuse LLM trace sync (incremental)
After=network.target hapax-secrets.service
Requires=hapax-secrets.service
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/home/hapax/projects/hapax-council
EnvironmentFile=/run/user/1000/hapax-secrets.env
ExecStart=/home/hapax/.local/bin/uv run python -m agents.langfuse_sync --auto
Environment=PATH=/home/hapax/.local/bin:/home/hapax/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
MemoryMax=512M
CPUQuota=50%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=langfuse-sync
TimeoutStartSec=1200
