# systemd/units-pi6/obsidian-sync.service
[Unit]
Description=Obsidian vault sync (Pi6 offload)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/hapax/bin/sync-agent-wrapper.sh obsidian_sync Documents/Personal/
MemoryMax=512M
CPUQuota=50%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=obsidian-sync

[Install]
WantedBy=default.target
