# systemd/overrides/studio-compositor.service.d/limit-nofile.conf
[Service]
# Drop #41 finding 5: FD leak causing EMFILE core dumps via GLib wakeup
# pipe creation. Compositor core-dumped three times in 4 minutes on
# 2026-04-14 with OSError: [Errno 24] Too many open files from routine
# writes (status.tmp, snapshot.jpg.tmp, last-ntfy-*.txt) and finally
# GLib-ERROR: Creating pipes for GWakeup: Too many open files. systemd
# user services default to 1024-4096 fds. 65536 gives 16x headroom
# while the underlying leak (unidentified — HLS sink internal fds or
# pipeline-rebuild cascade per drops #27/#37 most suspected) is tracked
# down via /proc/self/fd census.
LimitNOFILE=65536
