# systemd/overrides/studio-compositor.service.d/ordering.conf
[Unit]
# All services managed by systemd directly — no process-compose
After=hapax-daimonion.service visual-layer-aggregator.service hapax-secrets.service
Wants=visual-layer-aggregator.service
# Clear the original Wants=hapax-daimonion.service from base unit (voice starts via default.target)
Wants=
Wants=visual-layer-aggregator.service

[Service]
# Brief delay — camera USB enumeration can race with GStreamer pipeline init
ExecStartPre=/bin/sleep 10
TimeoutStartSec=120
