# systemd/overrides/studio-compositor.service.d/ordering.conf
[Unit]
# Voice is managed by process-compose (inside hapax-stack), not standalone systemd
After=hapax-stack.service
Wants=hapax-stack.service
# Clear the original Wants=hapax-voice.service that starts a duplicate daemon
Wants=
Wants=hapax-stack.service

[Service]
# Brief delay — camera USB enumeration can race with GStreamer pipeline init
ExecStartPre=/bin/sleep 10
TimeoutStartSec=120
