# systemd/overrides/studio-compositor.service.d/priority.conf
[Service]
# Tier 0: Real-time video path — highest CPU priority
CPUWeight=500
IOWeight=400
