# systemd/overrides/studio-compositor.service.d/ward-modulator.conf
[Service]
# Ward stimmung modulator activation (z-axis spec Phase 2 PR #1147).
# Runs at ~5 Hz inside fx_tick_callback, reads imagination dimensions
# from /dev/shm/hapax-imagination/current.json, and writes per-ward
# z_index_float + alpha for non-default-plane wards. Default-plane
# (``on-scrim``) wards remain untouched (director / recruitment owns).
#
# When unset, ``WardStimmungModulator.maybe_tick`` returns immediately
# so the compositor sees no behavior change. With ``=1`` the modulator
# starts ticking and emits hapax_ward_modulator_tick_total /
# hapax_ward_depth_attenuation metrics for Grafana to plot.
# 2026-04-23 operator-flagged CBIP + HARDM blinking at ~500ms sinewave.
# Root cause: modulator writes per-ward alpha from imagination depth every
# ~5 Hz with TTL 0.4s → `paint_with_alpha(props.alpha)` modulates every
# blit → visible blink on non-default-plane wards (HARDM beyond-scrim,
# album_overlay mid-scrim). Disabled via PR #1247 pending a hysteresis/
# epsilon-gate fix.
#
# 2026-04-23 re-enabled after the epsilon-gate landed in
# ``ward_stimmung_modulator.py``: WARD_PROPERTIES_TTL_S bumped 0.4 → 2.5s
# (covers ~12 tick cycles of scheduler slack so entries never expire
# between writes), MIN_DELTA=0.02 epsilon-gate stops micro-oscillation
# churn. With both in place, depth attenuation resumes producing the
# spec-§4 z-plane stratification without the 5 Hz alpha oscillation
# that caused the original blink. Reversible: flip back to =0.
Environment=HAPAX_WARD_MODULATOR_ACTIVE=1
