# systemd/overrides/studio-compositor.service.d/ward-modulator.conf
[Service]
# Ward stimmung modulator activation (z-axis spec Phase 2 PR #1147).
# Runs at ~5 Hz inside fx_tick_callback, reads imagination dimensions
# from /dev/shm/hapax-imagination/current.json, and writes per-ward
# z_index_float + alpha for non-default-plane wards. Default-plane
# (``on-scrim``) wards remain untouched (director / recruitment owns).
#
# When unset, ``WardStimmungModulator.maybe_tick`` returns immediately
# so the compositor sees no behavior change. With ``=1`` the modulator
# starts ticking and emits hapax_ward_modulator_tick_total /
# hapax_ward_depth_attenuation metrics for Grafana to plot.
Environment=HAPAX_WARD_MODULATOR_ACTIVE=1
