# systemd/overrides/rag-ingest.service.d/ordering.conf
[Unit]
After=network.target llm-stack.service llm-stack-analytics.service
Wants=llm-stack.service llm-stack-analytics.service
