# systemd/overrides/rag-ingest.service.d/priority.conf
[Service]
# Tier 2: Background batch work — yields to everything interactive
Nice=10
CPUWeight=25
IOWeight=50
