# systemd/hapax-cc-pr-merge-watcher.service
[Unit]
Description=Hapax CC-task PR-merge auto-close watcher (H9, PR3 of cc-hygiene)
Documentation=file://%h/projects/hapax-council/docs/research/2026-04-26-task-list-hygiene-operator-visibility.md
After=hapax-secrets.service

[Service]
Type=oneshot
ExecStart=%h/projects/hapax-council/.venv/bin/python %h/projects/hapax-council/scripts/cc-pr-merge-watcher.py
WorkingDirectory=%h/projects/hapax-council
Environment=HOME=%h
Environment=PYTHONPATH=%h/projects/hapax-council
# Killswitch: set HAPAX_CC_HYGIENE_OFF=1 in the environment to skip both
# the H8 PostToolUse hook and this watcher.

[Install]
WantedBy=default.target
