# systemd/hapax-rebuild-logos.timer
[Unit]
Description=Rebuild Hapax Logos/Reverie binaries (5 min check)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
