# systemd/units/audio-recorder.service
[Unit]
Description=Continuous ambient audio recording (Blue Yeti via PipeWire)
After=pipewire.service pipewire-pulse.service
Requires=pipewire-pulse.service

[Service]
Type=simple
ExecStartPre=/usr/bin/mkdir -p %h/audio-recording/raw
ExecStart=/usr/bin/ffmpeg -nostdin -f pulse \
  -i alsa_input.usb-Blue_Microphones_Yeti_Stereo_Microphone_REV8-00.analog-stereo \
  -ac 1 -f segment -segment_time 900 -strftime 1 -c:a flac \
  %h/audio-recording/raw/rec-%%Y%%m%%d-%%H%%M%%S.flac
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audio-recorder

[Install]
WantedBy=default.target
