# systemd/units/audio-recorder.service
[Unit]
Description=Continuous ambient audio recording (Blue Yeti via PipeWire)
After=pipewire.service pipewire-pulse.service
Requires=pipewire-pulse.service

[Service]
Type=simple
ExecStartPre=/usr/bin/mkdir -p %h/audio-recording/raw
ExecStart=/bin/bash -c 'exec /usr/bin/pw-record \
  --target alsa_input.usb-Blue_Microphones_Yeti_Stereo_Microphone_REV8-00.analog-stereo \
  --format s16 --rate 48000 --channels 1 - | \
  /usr/bin/ffmpeg -nostdin -f s16le -ar 48000 -ac 1 -i pipe: \
  -c:a flac -f segment -segment_time 900 -strftime 1 \
  %h/audio-recording/raw/rec-%%Y%%m%%d-%%H%%M%%S.flac'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audio-recorder

[Install]
WantedBy=default.target
