# systemd/units/claude-code-sync.timer
[Unit]
Description=Claude Code sync (every 6h)

[Timer]
OnBootSec=6min
OnUnitActiveSec=6h
RandomizedDelaySec=10min
Persistent=true

[Install]
WantedBy=timers.target
