# systemd/units/claude-md-audit.service
[Unit]
Description=Monthly CLAUDE.md rotation policy + sister-extension audit
Documentation=file:%h/projects/hapax-council/docs/superpowers/specs/2026-04-13-claude-md-excellence-design.md

[Service]
Type=oneshot
ExecStart=%h/projects/hapax-council/scripts/monthly-claude-md-audit.sh
StandardOutput=journal
StandardError=journal
Nice=10
