# systemd/units/claude-md-audit.timer
[Unit]
Description=Monthly CLAUDE.md audit schedule

[Timer]
OnCalendar=monthly
RandomizedDelaySec=1h
Persistent=true
Unit=claude-md-audit.service

[Install]
WantedBy=timers.target
