# systemd/units/daily-briefing.service
[Unit]
Description=Daily system briefing generator
After=network.target digest.service
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/home/hapax/projects/ai-agents
ExecStart=/home/hapax/.local/bin/briefing-watchdog
Environment=PATH=/home/hapax/.local/bin:/home/hapax/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=GNUPGHOME=/home/hapax/.gnupg
Environment=PASSWORD_STORE_DIR=/home/hapax/.password-store
MemoryMax=2G
CPUQuota=60%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daily-briefing
