# systemd/units/deliberation-eval.path
[Path]
PathChanged=/home/hapax/projects/hapax-council/profiles/deliberations
Unit=deliberation-eval.service

[Install]
WantedBy=default.target
