# systemd/units/dev-story-index.service
[Unit]
Description=Dev story — rebuild development archaeology index
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/home/hapax/projects/hapax-council
ExecStart=/home/hapax/.local/bin/uv run python -m agents.dev_story --index
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dev-story-index
MemoryMax=512M
CPUQuota=50%
