# systemd/units/digest.timer
[Unit]
Description=Daily content digest at 06:45

[Timer]
OnCalendar=*-*-* 06:45:00
Persistent=true

[Install]
WantedBy=timers.target
