# systemd/units/disk-space-check.service
[Unit]
Description=Disk space monitor — alert at 85%% and 95%%

[Service]
Type=oneshot
ExecStart=/home/hapax/projects/hapax-council/scripts/disk-space-check.sh
