# systemd/units/drift-detector.timer
[Unit]
Description=Twice-daily documentation drift reconciliation (03:00, 15:00)

[Timer]
OnCalendar=*-*-* 03:00:00
OnCalendar=*-*-* 15:00:00
RandomizedDelaySec=10min
Persistent=true

[Install]
WantedBy=timers.target
