# systemd/units/flow-journal.service
[Unit]
Description=Flow journal — persist flow state transitions to RAG
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/home/operator/projects/hapax-council
ExecStart=/home/operator/.local/bin/uv run python -m agents.flow_journal
Environment=PATH=/home/operator/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flow-journal
MemoryMax=128M
CPUQuota=20%
