# systemd/units/git-sync.service
[Unit]
Description=Git commit history RAG sync (incremental)
After=network.target
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
ExecStart=/home/operator/.local/bin/uv run --project /home/operator/projects/hapax-council python -m agents.git_sync --auto
Environment=PATH=/home/operator/.local/bin:/home/operator/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=GNUPGHOME=/home/operator/.gnupg
Environment=PASSWORD_STORE_DIR=/home/operator/.password-store
MemoryMax=256M
CPUQuota=50%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=git-sync
