# systemd/units/gpg-keyboxd-watchdog.timer
[Unit]
Description=Check keyboxd health every 5 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=300
AccuracySec=30

[Install]
WantedBy=timers.target
