# systemd/units/hapax-dmn.service
[Unit]
Description=Hapax DMN — always-on cognitive substrate
After=hapax-secrets.service ollama.service
Wants=ollama.service
PartOf=hapax-visual-stack.target

[Service]
Type=simple
ExecStart=/home/hapax/.local/bin/uv run python -m agents.dmn
WorkingDirectory=/home/hapax/projects/hapax-council
Restart=on-failure
RestartSec=10
EnvironmentFile=/run/user/1000/hapax-secrets.env
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin
MemoryMax=1G

[Install]
WantedBy=default.target
