# systemd/units/hapax-lufs-panic-cap.service
[Unit]
Description=Hapax broadcast master LUFS-S panic-cap (defense-in-depth)
Documentation=file:///home/hapax/projects/hapax-council/agents/studio_compositor/lufs_panic_cap.py
After=pipewire.service wireplumber.service hapax-secrets.service
Requires=pipewire.service wireplumber.service
PartOf=studio-compositor.service

[Service]
Type=simple
WorkingDirectory=/home/hapax/projects/hapax-council
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/home/hapax/projects/hapax-council/.venv/bin/python -m agents.studio_compositor.lufs_panic_cap
Restart=always
RestartSec=5
# Resource caps: numpy + parec is light; reserve modest memory and a
# single CPU. Prometheus HTTP server on :9484 (pinned in lufs_panic_cap.py).
MemoryMax=512M
CPUWeight=50
# Crash visibility — let journal capture both streams.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
