# systemd/units/hapax-discord-webhook.service
[Unit]
Description=Hapax Discord webhook poster (ytb-010 Phase 1, cross-surface federation)
Documentation=https://github.com/ryanklee/hapax-council/blob/main/agents/cross_surface/discord_webhook.py
After=network-online.target hapax-secrets.service
Wants=network-online.target
PartOf=hapax.target

[Service]
Type=simple
WorkingDirectory=%h/projects/hapax-council
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=PATH=%h/projects/hapax-council/.venv/bin:/usr/bin
Environment=HAPAX_DISCORD_WEBHOOK_METRICS_PORT=9500
Environment=HAPAX_DISCORD_WEBHOOK_TICK_S=30
Environment=HAPAX_LOG_LEVEL=INFO
# Operator must export HAPAX_DISCORD_WEBHOOK_URL via hapax-secrets:
#   pass insert hapax/discord-webhook-url   # then echo into env file
# Without it, the daemon idles + logs "no_webhook_url" per event.
ExecStart=%h/projects/hapax-council/.venv/bin/python -m agents.cross_surface
Restart=on-failure
RestartSec=15
MemoryMax=128M

[Install]
WantedBy=default.target
