# systemd/units/hapax-video-cam@.service
[Unit]
Description=Video capture — %i camera (ffmpeg segmented MKV)
After=network.target
OnFailure=notify-failure@%n.service

[Service]
Type=simple
EnvironmentFile=%h/.config/hapax-video/%i.env
ExecStartPre=/usr/bin/mkdir -p %h/video-recording/%i
ExecStart=/home/hapax/.local/bin/uv run --directory /home/hapax/projects/hapax-council python -m agents.video_capture --camera %i --device ${DEVICE}
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Restart=on-failure
RestartSec=10
MemoryMax=512M
CPUQuota=30%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=video-cam-%i

[Install]
WantedBy=default.target
