# systemd/units/hapax-secrets.service
[Unit]
Description=Hapax Secrets — load credentials from pass store
# No dependencies — this is the first thing that runs

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/bash -c 'install -m 600 /dev/null /run/user/1000/hapax-secrets.env'
ExecStart=/bin/bash -c '\
  echo "LITELLM_API_KEY=$(pass show litellm/master-key)" >> /run/user/1000/hapax-secrets.env && \
  echo "LANGFUSE_PUBLIC_KEY=$(pass show langfuse/public-key)" >> /run/user/1000/hapax-secrets.env && \
  echo "LANGFUSE_SECRET_KEY=$(pass show langfuse/secret-key)" >> /run/user/1000/hapax-secrets.env && \
  echo "HF_TOKEN=$(pass show api/huggingface)" >> /run/user/1000/hapax-secrets.env && \
  echo "LITELLM_BASE_URL=http://127.0.0.1:4000" >> /run/user/1000/hapax-secrets.env && \
  echo "LANGFUSE_HOST=http://127.0.0.1:3000" >> /run/user/1000/hapax-secrets.env && \
  echo "MISTRAL_API_KEY=$(pass show api/mistral)" >> /run/user/1000/hapax-secrets.env \
'
ExecStopPost=/bin/rm -f /run/user/1000/hapax-secrets.env
Environment=GNUPGHOME=/home/hapax/.gnupg
Environment=PASSWORD_STORE_DIR=/home/hapax/.password-store
Environment=HOME=/home/hapax
SyslogIdentifier=hapax-secrets

[Install]
WantedBy=default.target
