# systemd/units/hapax-hardm-publisher.service
[Unit]
Description=HARDM signal publisher — writes 16 primary signals to SHM for the dot-matrix ward
After=hapax-secrets.service studio-compositor.service
Wants=studio-compositor.service

[Service]
Type=oneshot
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=HARDM_PUBLISH_LOG_LEVEL=WARNING
ExecStart=%h/.local/bin/uv run python %h/projects/hapax-council/scripts/hardm-publish-signals.py
WorkingDirectory=%h/projects/hapax-council
TimeoutSec=10
MemoryMax=128M
