# systemd/units/hapax-audio-router.service
[Unit]
Description=Hapax dynamic audio router (Evil Pet + S-4 5 Hz arbiter)
Documentation=https://github.com/ryanklee/hapax-council/blob/main/docs/superpowers/specs/2026-04-21-evilpet-s4-dynamic-dual-processor-design.md
After=hapax-secrets.service pipewire.service
Wants=hapax-secrets.service

[Service]
Type=simple
WorkingDirectory=%h/projects/hapax-council
EnvironmentFile=-%h/.cache/hapax/secrets.env
ExecStart=%h/projects/hapax-council/.venv/bin/python -m agents.audio_router.dynamic_router
Restart=on-failure
RestartSec=2s
# Audio router is small and pure-Python — generous bound vs the
# kernel-tuning baseline (vm.min_free_kbytes 512MB) so a single tick
# stays well under the global OOM watermark.
MemoryHigh=512M
MemoryMax=1G

# Pure software; no GPU. Isolate from CUDA so a stray import cannot
# acquire VRAM the compositor / TabbyAPI rely on.
Environment=CUDA_VISIBLE_DEVICES=

# Stdout / stderr → systemd journal (default).
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
