# systemd/units/hapax-daimonion.service
[Unit]
Description=Hapax Daimonion — persistent voice interaction daemon
After=pipewire.service pipewire-pulse.service hapax-secrets.service
Requires=pipewire.service hapax-secrets.service
Wants=pipewire-pulse.service
StartLimitBurst=8
StartLimitIntervalSec=600
OnFailure=notify-failure@%n.service

[Service]
Type=simple
WorkingDirectory=/home/hapax/projects/hapax-council
# Use .venv/bin/python directly instead of `uv run` because NeMo toolkit
# requires lightning==2.4.0 (nv-one-logger crash with >=2.5), and `uv run`
# re-syncs to pyproject.toml which installs lightning 2.6.x.
# NeMo is installed via `uv pip install` outside pyproject.toml resolution
# due to huggingface-hub version conflict with pydantic-ai.
# Delay startup to let logos-api initialize first,
# preventing GPU memory spike that causes SIGKILL during boot
ExecStartPre=/bin/sleep 10
ExecStart=/home/hapax/projects/hapax-council/.venv/bin/python -m agents.hapax_daimonion
Restart=always
RestartSec=10
EnvironmentFile=/run/user/1000/hapax-secrets.env
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=OMP_NUM_THREADS=2
Environment=MKL_NUM_THREADS=2
Environment=OPENBLAS_NUM_THREADS=2
Environment=ONNXRUNTIME_INTRA_OP_NUM_THREADS=2
MemoryMax=8G
OOMScoreAdjust=-500
Nice=-10
LimitNOFILE=1048576
SyslogIdentifier=hapax-daimonion

[Install]
WantedBy=default.target
