# systemd/units/hapax-refused-lifecycle-constitutional.service
[Unit]
Description=hapax-refused-lifecycle-constitutional — type-B inotify watcher on axiom registry / manifesto / refusal-brief / MEMORY.md
Documentation=https://github.com/ryanklee/hapax-council/blob/main/agents/refused_lifecycle/constitutional_watcher.py

# Always-on filesystem-event-driven daemon. NO clock-poll — axioms are
# stable until amended; inotify is the precisely-tractable trigger.
# On any change to a watched constitutional surface, walks all type-B
# REFUSED cc-tasks and re-evaluates against the current axiom state.
#
# Probe semantics (typical type-B): lift_polarity=absent — keyword
# *removed* from the file signals lift. The cold-contact-alphaxiv-comments
# multi-classified slug uses present-polarity for its constitutional half.

[Service]
Type=simple
WorkingDirectory=%h/projects/hapax-council
EnvironmentFile=-%h/.config/hapax/env
ExecStart=%h/projects/hapax-council/.venv/bin/python -m agents.refused_lifecycle.constitutional_watcher
Restart=always
RestartSec=10s

[Install]
WantedBy=default.target
