# systemd/units/hapax-youtube-viewer-count.service
[Unit]
Description=Hapax YouTube viewer-count producer (orphan-ward Phase 4)
Documentation=docs/superpowers/plans/2026-04-20-orphan-ward-producers-plan.md
After=hapax-secrets.service
Requires=hapax-secrets.service

[Service]
Type=simple
Restart=on-failure
RestartSec=10
TimeoutStopSec=5
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=HOME=%h
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
ExecStart=/home/hapax/projects/hapax-council/.venv/bin/python /home/hapax/projects/hapax-council/scripts/youtube-viewer-count-producer.py

[Install]
WantedBy=default.target
