# systemd/units/hapax-environmental-emphasis.service
[Unit]
Description=Hapax Environmental Emphasis — salience-driven hero-mode driver
After=hapax-secrets.service studio-compositor.service
Wants=studio-compositor.service

[Service]
Type=oneshot
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=HAPAX_LOG_LEVEL=INFO
ExecStart=%h/.local/bin/uv run python %h/projects/hapax-council/scripts/environmental_emphasis_tick.py
WorkingDirectory=%h/projects/hapax-council
MemoryMax=256M
