# systemd/units/hapax-stream-auto-private.timer
[Unit]
Description=15s cadence auto-private enforcer (LRR Phase 6 §5)

[Timer]
OnBootSec=30s
OnUnitActiveSec=15s
Persistent=false

[Install]
WantedBy=timers.target
