# systemd/units/hapax-pin-check.service
[Unit]
Description=Hapax Ryzen pin-glitch live probe — one-shot detection + auto-fix
Documentation=file://%h/projects/hapax-council/scripts/hapax-pin-check-probe.sh
# Don't fire while pipewire is restarting; the pactl probe would
# misread sink state during the brief restart window.
After=pipewire.service pipewire-pulse.service
PartOf=pipewire.service

[Service]
Type=oneshot
# `uv run python` (the CLI's shebang) discovers the project + venv
# from the current working directory. Without WorkingDirectory the
# service runs from / and the CLI fails to import shared.* modules.
WorkingDirectory=%h/projects/hapax-council
# Wrapper assembles the live probe from pactl + pw-cat then invokes
# `hapax-audio-topology pin-check` with persisted state. Exit code 0
# = healthy or recovery succeeded; exit 1 = PIN_GLITCH detected and
# auto-fix declined or failed.
ExecStart=/home/hapax/projects/hapax-council/scripts/hapax-pin-check-probe.sh
# Allow the auto-fix path to invoke pactl set-card-profile twice (off
# then on); brief but non-trivial.
TimeoutStartSec=15s
# Don't restart on fail — the timer drives the cadence; a flapping
# detector should surface via journal grep, not an exec storm.
Restart=no

[Install]
WantedBy=default.target
