# systemd/units/hapax-watch-receiver.service
[Unit]
Description=Hapax Watch Sensor Receiver
After=network.target
OnFailure=notify-failure@%n.service

[Service]
Type=simple
WorkingDirectory=/home/hapax/projects/hapax-council
ExecStart=/home/hapax/projects/hapax-council/.venv/bin/uvicorn agents.watch_receiver:app --host 0.0.0.0 --port 8042
Environment=HOME=/home/hapax
Environment=HAPAX_HOME=/home/hapax
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hapax-watch-receiver
MemoryMax=256M

[Install]
WantedBy=default.target
