# systemd/units/hapax-rebuild-services.service
[Unit]
Description=Hapax Python services — rebuild from main
After=network-online.target

[Service]
Type=oneshot
Environment=HOME=%h
Environment=PATH=%h/.local/bin:%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin
SyslogIdentifier=hapax-rebuild-services
TimeoutStartSec=120

# hapax-dmn: restart if DMN, imagination, or reverie code changed
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-council --service hapax-dmn.service --watch "agents/dmn/ agents/imagination.py agents/imagination_resolver.py agents/imagination_context.py agents/reverie/ agents/visual_chain.py agents/effect_graph/ shared/" --sha-key dmn

# hapax-daimonion: restart if voice agent or shared code changed
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-council --service hapax-daimonion.service --watch "agents/hapax_daimonion/ shared/" --sha-key voice

# logos-api: restart if logos backend code changed
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-council --service logos-api.service --watch "logos/" --sha-key logos-api

# officium-api: restart if officium code changed (separate repo)
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-officium --service officium-api.service --sha-key officium

# hapax-mcp: pull-only (MCP server restarts on Claude Code reconnect)
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-mcp --sha-key hapax-mcp --pull-only

# studio-compositor: restart if compositor, effect-graph, shader nodes, or
# preset code changed. BETA-FINDING-P (queue 025 Phase 6) — migrates the
# compositor off the local ``studio-compositor-reload.path`` unit, which
# fired on any file-system change without a branch check and caused three
# unexplained restarts during queue 023 research while alpha was editing
# the compositor on a feature branch. rebuild-service.sh runs the
# branch-check and skips the deploy on feature branches instead.
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-council --service studio-compositor.service --watch "agents/studio_compositor/ agents/effect_graph/ shared/" --sha-key compositor

# ── Gap 2 (round-4 Phase 6 coverage expansion): 8 daemons that previously
# had no rebuild coverage. The editing-then-restart pattern was manual for
# all of them; now they are picked up by the 5-minute
# hapax-rebuild-services.timer with the same branch-check discipline as
# the services above. Watch lists are minimised to the subdirectories
# each daemon actually imports so an unrelated edit to another module
# under ``shared/`` still triggers every service — that is intentional
# because the shared layer is cross-cutting.

# visual-layer-aggregator: perception fusion feeding stimmung
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-council --service visual-layer-aggregator.service --watch "agents/visual_layer_aggregator/ shared/" --sha-key vla

# hapax-content-resolver: content scheduling / ground-surface enrichment
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-council --service hapax-content-resolver.service --watch "agents/content_resolver/ shared/" --sha-key content-resolver

# hapax-watch-receiver: biometric ingestion (watch_receiver is a single module, not a package)
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-council --service hapax-watch-receiver.service --watch "agents/watch_receiver.py shared/" --sha-key watch-receiver

# hapax-reverie: visual expression Python daemon (distinct from the Rust wgpu binary)
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-council --service hapax-reverie.service --watch "agents/reverie/ shared/" --sha-key reverie

# hapax-imagination-loop: imagination daemon loop (single-file modules)
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-council --service hapax-imagination-loop.service --watch "agents/imagination.py agents/imagination_loop.py agents/imagination_resolver.py agents/imagination_context.py shared/" --sha-key imagination-loop

# studio-fx retired 2026-04-17: the numpy/OpenCV effects pipeline was
# orphaned after the effect_graph/WGSL unification. Every fx-*.jpg
# output had zero readers, and the virtualcam target (/dev/video50) is
# owned by the YouTube slot-0 ffmpeg. Freed ~1 core of wasted CPU.
# The effect_graph watch group (agents/effect_graph/) is still covered
# via studio-compositor's own rebuild entry above.

# studio-fx-output: shell-launcher driving the GStreamer output bin; the
# ExecStart IS the shell script, so its own source is the watch path.
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-council --service studio-fx-output.service --watch "systemd/units/studio-fx-output.sh" --sha-key studio-fx-output

# studio-person-detector: camera-side person detection
# (agents/studio_person_detector.py is the single-file module)
ExecStart=%h/projects/hapax-council/scripts/rebuild-service.sh --repo %h/projects/hapax-council --service studio-person-detector.service --watch "agents/studio_person_detector.py shared/" --sha-key studio-person-detector
