# systemd/units/hapax-imagination-loop.service
[Unit]
Description=Hapax Imagination Loop Daemon
After=hapax-secrets.service hapax-dmn.service
Wants=hapax-dmn.service
PartOf=hapax-visual-stack.target

[Service]
Type=simple
EnvironmentFile=/run/user/1000/hapax-secrets.env
Environment=OMP_NUM_THREADS=2
Environment=MKL_NUM_THREADS=2
ExecStart=%h/.local/bin/uv run python -m agents.imagination_daemon
WorkingDirectory=%h/projects/hapax-council
MemoryMax=1G
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300
OnFailure=notify-failure@%n.service

[Install]
WantedBy=hapax-visual-stack.target
