# systemd/units/hapax-egress-audit-rotate.service
[Unit]
Description=hapax — daily rotate + prune of demonet egress audit JSONL
Documentation=file:///home/hapax/projects/hapax-council/docs/superpowers/plans/2026-04-20-demonetization-safety-plan.md
OnFailure=notify-failure@%n.service
ExecStart=/home/hapax/projects/hapax-council/scripts/hapax-egress-audit-rotate.sh

[Service]
Type=oneshot
Nice=19
StandardOutput=journal
StandardError=journal
