# systemd/units/hapax-hardm-publisher.timer
[Unit]
Description=Publish HARDM cell signals to SHM every 2 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=2s
Unit=hapax-hardm-publisher.service
AccuracySec=500ms

[Install]
WantedBy=timers.target
