# systemd/units/hapax-environmental-emphasis.timer
[Unit]
Description=Tick the environmental-salience emphasis driver every 30s

[Timer]
OnBootSec=60s
OnUnitActiveSec=30s
Unit=hapax-environmental-emphasis.service
AccuracySec=5s

[Install]
WantedBy=timers.target
