# systemd/units/hapax-stack.service
[Unit]
Description=Hapax Application Stack (process-compose)
After=docker.service pipewire.service logos-api.service

[Service]
Type=simple
WorkingDirectory=/home/hapax/projects/hapax-council
ExecStart=/usr/bin/process-compose up --tui=false --port 9080
ExecStop=/usr/bin/process-compose down --port 9080
Restart=on-failure
RestartSec=10
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=GNUPGHOME=/home/hapax/.gnupg
Environment=PASSWORD_STORE_DIR=/home/hapax/.password-store

[Install]
WantedBy=default.target
