# systemd/units/hapax-bluesky-post.service
[Unit]
Description=Hapax Bluesky poster (ytb-010 Phase 2, cross-surface federation)
Documentation=https://github.com/ryanklee/council/blob/main/agents/cross_surface/bluesky_post.py
After=network-online.target hapax-secrets.service
Wants=network-online.target
PartOf=hapax.target

[Service]
Type=simple
WorkingDirectory=%h/projects/hapax-council
EnvironmentFile=-/run/user/1000/hapax-secrets.env
Environment=PATH=%h/projects/hapax-council/.venv/bin:/usr/bin
Environment=HAPAX_BLUESKY_METRICS_PORT=9501
Environment=HAPAX_BLUESKY_TICK_S=30
Environment=HAPAX_LOG_LEVEL=INFO
# Operator must export from hapax-secrets:
#   HAPAX_BLUESKY_HANDLE         (e.g. hapax.bsky.social)
#   HAPAX_BLUESKY_APP_PASSWORD   (19-char from https://bsky.app/settings/app-passwords)
# Without either, daemon idles + logs "no_credentials" per event.
ExecStart=%h/projects/hapax-council/.venv/bin/python -m agents.cross_surface.bluesky_post
Restart=on-failure
RestartSec=15
MemoryMax=128M

[Install]
WantedBy=default.target
