# systemd/units/hapax-logos.service
[Unit]
Description=Hapax Logos — native visual + wgpu rendering surface
After=graphical-session.target logos-api.service hapax-imagination.service
Wants=graphical-session.target hapax-imagination.service
PartOf=hapax-visual-stack.target
# Replaces standalone hapax-visual (visual surface is now a second window)
Conflicts=hapax-visual.service
StartLimitBurst=5
StartLimitIntervalSec=300
OnFailure=notify-failure@%n.service

[Service]
Type=simple
ExecStartPre=/usr/bin/mkdir -p /dev/shm/hapax-logos /dev/shm/hapax-visual /dev/shm/hapax-compositor /dev/shm/hapax-stimmung
ExecStartPre=/home/hapax/projects/hapax-council/scripts/bootstrap-profiles.sh
ExecStart=/home/hapax/.local/bin/hapax-logos
Restart=on-failure
RestartSec=5
SuccessExitStatus=SIGTERM
Environment=RUST_LOG=info
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=HOME=/home/hapax
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
# WebKit2GTK needs dbus
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
# Vulkan for wgpu visual surface
Environment=VK_ICD_FILENAMES=/usr/share/vulkan/icd.d/nvidia_icd.json
# NVIDIA explicit sync triggers webkit2gtk protocol error on Wayland (gtk#8056, tauri#10702)
# __NV_DISABLE_EXPLICIT_SYNC keeps GPU accel; WEBKIT_DISABLE_DMABUF_RENDERER is software fallback
# See docs/issues/tauri-wayland-protocol-error.md
Environment=__NV_DISABLE_EXPLICIT_SYNC=1
# Disable DMA-BUF renderer — native Wayland with NVIDIA triggers SIGKILL via
# webkit2gtk syncobj bug (gtk#8056). GDK_BACKEND=x11 avoids the crash but GBM
# buffer creation fails under XWayland, producing a black webview.
# WEBKIT_DISABLE_DMABUF_RENDERER forces software compositing for the webview
# while keeping native Wayland for window management.
Environment=WEBKIT_DISABLE_DMABUF_RENDERER=1
MemoryMax=4G
SyslogIdentifier=hapax-logos

[Install]
WantedBy=graphical-session.target
