# systemd/units/hapax-reverie.service
[Unit]
Description=Hapax Reverie — visual expression daemon
After=hapax-secrets.service hapax-dmn.service
Wants=hapax-dmn.service
PartOf=hapax-visual-stack.target

[Service]
Type=simple
ExecStart=/home/hapax/.local/bin/uv run python -m agents.reverie
WorkingDirectory=/home/hapax/projects/hapax-council
Restart=on-failure
RestartSec=10
EnvironmentFile=/run/user/1000/hapax-secrets.env
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin
MemoryMax=1G

[Install]
WantedBy=default.target
