# systemd/units/health-monitor.service
[Unit]
Description=LLM Stack Health Monitor (auto-fix + notify)
After=network.target
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
SuccessExitStatus=1
ExecStart=/home/hapaxlegomenon/.local/bin/health-watchdog
Environment=PATH=/home/hapaxlegomenon/.local/bin:/home/hapaxlegomenon/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapaxlegomenon
Environment=GNUPGHOME=/home/hapaxlegomenon/.gnupg
Environment=PASSWORD_STORE_DIR=/home/hapaxlegomenon/.password-store
StandardOutput=journal
StandardError=journal
SyslogIdentifier=health-monitor
MemoryMax=512M
CPUQuota=30%
