# systemd/units/health-monitor.timer
[Unit]
Description=LLM Stack Health Monitor (every 15 min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
RandomizedDelaySec=1min
Persistent=true

[Install]
WantedBy=timers.target
