# systemd/units/health-connect-parse.timer
[Unit]
Description=Daily Health Connect parsing

[Timer]
OnBootSec=30min
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target
