# systemd/units/health-monitor.service
[Unit]
Description=LLM Stack Health Monitor (auto-fix + notify)
After=network.target
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
SuccessExitStatus=1 2
ExecStart=/home/operator/.local/bin/health-watchdog
Environment=PATH=/home/operator/.local/bin:/home/operator/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=GNUPGHOME=/home/operator/.gnupg
Environment=PASSWORD_STORE_DIR=/home/operator/.password-store
StandardOutput=journal
StandardError=journal
SyslogIdentifier=health-monitor
TimeoutStartSec=300
MemoryMax=768M
CPUQuota=50%
