# systemd/units/hls-archive-rotate.service
[Unit]
Description=LRR Phase 2 HLS archive rotation pass
After=studio-compositor.service
Wants=studio-compositor.service

[Service]
Type=oneshot
WorkingDirectory=%h/projects/hapax-council
ExecStart=%h/projects/hapax-council/.venv/bin/python %h/projects/hapax-council/scripts/hls-archive-rotate.py --json
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
TimeoutStartSec=120

[Install]
WantedBy=default.target
