# systemd/units/knowledge-maint.service
[Unit]
Description=Knowledge base maintenance — prune and deduplicate
After=network.target
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/home/hapax/projects/ai-agents
ExecStart=/home/hapax/.local/bin/knowledge-maint-watchdog
Environment=PATH=/home/hapax/.local/bin:/home/hapax/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=GNUPGHOME=/home/hapax/.gnupg
Environment=PASSWORD_STORE_DIR=/home/hapax/.password-store
MemoryMax=1G
CPUQuota=30%
StandardOutput=journal
StandardError=journal
SyslogIdentifier=knowledge-maint
