# systemd/units/llm-stack.service
[Unit]
Description=LLM Stack (Docker Compose full profile)
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/hapax/llm-stack
Environment=HOME=/home/hapax
Environment=GNUPGHOME=/home/hapax/.gnupg
Environment=PASSWORD_STORE_DIR=/home/hapax/.password-store
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/usr/bin/docker info -f '{{.ServerVersion}}'
ExecStartPre=/home/hapax/llm-stack/generate-env.sh
ExecStart=/usr/bin/docker compose --profile full up -d
ExecStop=/usr/bin/docker compose --profile full down
ExecStopPost=/bin/rm -f /home/hapax/llm-stack/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llm-stack
TimeoutStartSec=120

[Install]
WantedBy=default.target
