# systemd/units/log-anomaly-alert.service
[Unit]
Description=Log anomaly detection — alert on ERROR spikes per agent
After=default.target

[Service]
Type=oneshot
ExecStart=%h/projects/distro-work/scripts/log-anomaly-alert.sh
Environment=NTFY_URL=http://localhost:8090/hapax-alerts
