# systemd/units/log-anomaly-alert.timer
[Unit]
Description=Run log anomaly detection every 30 minutes

[Timer]
OnCalendar=*:00/30
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
