# systemd/units/mediamtx.service
[Unit]
Description=MediaMTX RTMP relay for studio livestream → YouTube
Documentation=https://mediamtx.org/
After=network-online.target hapax-secrets.service
Wants=network-online.target
PartOf=studio-compositor.service

[Service]
Type=simple
WorkingDirectory=%h/projects/hapax-council
ExecStart=%h/projects/hapax-council/scripts/mediamtx-start.sh
Restart=on-failure
RestartSec=5s
StartLimitBurst=10
StartLimitIntervalSec=600
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=%h
Environment=XDG_RUNTIME_DIR=/run/user/1000
SyslogIdentifier=mediamtx
OnFailure=notify-failure@%n.service

[Install]
WantedBy=default.target
