# systemd/units/midi-route.service
[Unit]
Description=ALSA MIDI Port Routing (MCP → Hardware)
After=sound.target
ConditionPathExists=/proc/asound/card0

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/hapax/.local/bin/midi-route
ExecStop=/usr/bin/aconnect -x

# Retry if hardware not yet enumerated
Restart=on-failure
RestartSec=5
StartLimitBurst=10

[Install]
WantedBy=default.target
