# systemd/units/rag-ingest.service
[Unit]
Description=RAG Document Ingestion Watchdog
After=network.target
OnFailure=notify-failure@%n.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/hapaxlegomenon/projects/ai-agents
Environment=PATH=/home/hapaxlegomenon/.local/bin:/usr/local/bin:/usr/bin:/bin
# Uses .venv-ingest (isolated from main .venv due to docling/pydantic-ai
# huggingface-hub version conflict). Create with: make setup-ingest-venv
ExecStart=/home/hapaxlegomenon/projects/ai-agents/.venv-ingest/bin/python -m agents.ingest
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rag-ingest
MemoryMax=4G
CPUQuota=80%

[Install]
WantedBy=default.target
