# systemd/units/rclone-gdrive-drop.timer
[Unit]
Description=Sync gdrive:drop every 5s

[Timer]
OnBootSec=5s
OnUnitActiveSec=5s
AccuracySec=1s

[Install]
WantedBy=timers.target
