# systemd/units/screen-context.service
[Unit]
Description=Screen context — capture active window for behavioral profiling
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/home/operator/projects/hapax-council
ExecStart=/home/operator/.local/bin/uv run python -m agents.screen_context
Environment=PATH=/home/operator/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/1000
StandardOutput=journal
StandardError=journal
SyslogIdentifier=screen-context
MemoryMax=128M
CPUQuota=20%
