# systemd/units/stack-maintenance.service
[Unit]
Description=LLM Stack Weekly Maintenance
After=docker.service

[Service]
Type=oneshot
ExecStart=%h/llm-stack/stack-maintenance.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stack-maintenance
