# systemd/units/stimmung-sync.timer
[Unit]
Description=Stimmung sync (every 15min)

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
