# systemd/units/stimmung-sync.service
[Unit]
Description=Stimmung sync — persist system self-state to RAG
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/home/operator/projects/hapax-council
ExecStart=/home/operator/.local/bin/uv run python -m agents.stimmung_sync
Environment=PATH=/home/operator/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stimmung-sync
MemoryMax=128M
CPUQuota=20%
