# systemd/units/studio-compositor.service
[Unit]
Description=Studio Compositor — GPU-accelerated camera tiling, recording, HLS streaming, and visual effects
After=hapax-voice.service
Wants=hapax-voice.service
Conflicts=hapax-video-cam@brio.service hapax-video-cam@c920-hardware.service hapax-video-cam@c920-room.service hapax-video-cam@c920-aux.service
StartLimitBurst=5
StartLimitIntervalSec=300
OnFailure=notify-failure@%n.service

[Service]
Type=simple
WorkingDirectory=/home/hapax/projects/hapax-council
ExecStartPre=/usr/bin/mkdir -p %h/video-recording %h/.cache/hapax-compositor/hls
ExecStartPre=/usr/bin/find %h/.cache/hapax-compositor/hls -type f -delete
ExecStartPre=/home/hapax/projects/hapax-council/systemd/units/studio-camera-setup.sh
ExecStart=/home/hapax/projects/hapax-council/.venv/bin/python -m agents.studio_compositor
Restart=on-failure
RestartSec=10
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
Environment=XDG_RUNTIME_DIR=/run/user/1000
MemoryMax=4G
SyslogIdentifier=studio-compositor

[Install]
WantedBy=default.target
