# systemd/units/studio-compositor.service
[Unit]
Description=Studio Compositor — GPU-accelerated camera tiling, recording, HLS streaming, and visual effects
After=hapax-daimonion.service
Wants=hapax-daimonion.service
Conflicts=hapax-video-cam@brio.service hapax-video-cam@c920-desk.service hapax-video-cam@c920-room.service hapax-video-cam@c920-overhead.service
# Raised from 5/300 — the old limit was exhausted in under 3 minutes when
# reconnect loops went rampant. 20 restarts per hour is enough for genuine
# infrastructure churn without running away on a stuck camera.
StartLimitBurst=20
StartLimitIntervalSec=3600
OnFailure=notify-failure@%n.service

[Service]
# Type=notify + WatchdogSec: the compositor sends WATCHDOG=1 every 20s gated
# on "at least one camera has frames within the last 20s". Liveness means
# frames are flowing, not just "process is running". See camera epic Phase 1.
Type=notify
NotifyAccess=main
WatchdogSec=60s
WorkingDirectory=%h/projects/hapax-council
ExecStartPre=/usr/bin/mkdir -p %h/video-recording %h/.cache/hapax-compositor/hls
ExecStartPre=/usr/bin/find %h/.cache/hapax-compositor/hls -type f -delete
ExecStartPre=%h/projects/hapax-council/scripts/studio-install-udev-rules.sh
ExecStartPre=%h/projects/hapax-council/systemd/units/studio-camera-setup.sh
# Phase 2 archive precheck — verifies hls-archive-rotate.timer is
# enabled + active and ensures the archive destination dirs exist.
# Non-fatal: a broken precheck WARNs to journald but never blocks
# compositor startup (livestream uptime > archive completeness).
# Delta 2026-04-14-lrr-phase-2-hls-archive-dormant.md identified that
# Phase 2 shipped code + unit files without install-units.sh being
# run — every segment was being deleted at hlssink2 max_files=15 for
# hours of runtime with zero observable signal.
ExecStartPre=%h/projects/hapax-council/scripts/studio-compositor-archive-precheck.sh
ExecStart=%h/projects/hapax-council/.venv/bin/python -m agents.studio_compositor
# Livestream-performance-map Sprint 6 F7 / W1.5: postmortem capture.
# Fires on every watchdog timeout, segfault, restart, or manual stop.
# Saves nvidia-smi dump + journal tail + /dev/shm state + last metrics
# scrape to ~/hapax-state/postmortem/studio-compositor/{ts}/. Never
# exits non-zero — a broken hook must not block the shutdown cascade.
ExecStopPost=%h/projects/hapax-council/scripts/studio-compositor-postmortem.sh
Restart=on-failure
RestartSec=10
SuccessExitStatus=SIGTERM
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=%h
Environment=XDG_RUNTIME_DIR=/run/user/1000
EnvironmentFile=/run/user/1000/hapax-secrets.env
# Phase 10 observability polish — lock the compositor's CUDA-visible
# device to index 0 so nvh264enc and cudacompositor can never drift to
# a different GPU across reboots. Matches delta's sprint-5-delta-audit
# C2+C3 findings. Defaults are correct today; this makes the intent
# durable in the unit file itself.
Environment=CUDA_VISIBLE_DEVICES=0
# Cap OpenTelemetry BatchSpanProcessor so Langfuse outages can't fill the
# Python heap with pending spans — same treatment as hapax-daimonion got
# in PR #731's OOM-cascade containment. Defaults are 2048 queue and 30 s
# export timeout; at studio's per-frame span rate a Langfuse stall can
# backfill the queue faster than the WatchdogSec=60 s tick.
Environment=OTEL_BSP_MAX_QUEUE_SIZE=256
Environment=OTEL_BSP_MAX_EXPORT_BATCH_SIZE=64
Environment=OTEL_BSP_EXPORT_TIMEOUT=2000
Environment=OTEL_BSP_SCHEDULE_DELAY=1000
# Working set under 6 live cameras + 24-slot GL shader chain + MJPEG
# decoders + three Cairo overlay threads + interpipe producer+fallback
# pipelines plateaus at ~3.7 GB steady-state and spikes to ~3.9 GB
# during state-machine recovery rebuilds. The old hard 4 GB wall left
# zero headroom — a single FSM recovery swap or a Langfuse-outage span
# burst (see OTEL tunings above) was within one allocation of the
# cgroup kill. Mirror the daimonion pattern from PR #731: raise the
# hard ceiling and add a MemoryHigh soft ceiling so the kernel reclaims
# aggressively before MemoryMax turns into a kill.
# Live-observed 2026-04-13 15:xx CDT: PR #741's ``MemoryHigh=5G`` pushed
# 16.6 GB of anon pages to swap (vs ~5.7 GB RSS) before the WatchdogSec
# heartbeat missed and systemd killed the unit. The local drop-in at
# ``~/.config/systemd/user/studio-compositor.service.d/memory-override.conf``
# that set ``MemoryHigh=infinity`` was keeping the live system stable,
# but that state lived only on one workstation. Promote the drop-in
# value into the repo so every redeploy picks it up without the
# operator having to re-install the override. Follow-up: the underlying
# ~360 MB/min anon leak is still open; track via
# ``studio_compositor_memory_footprint_bytes`` and the PR #770 RSS
# gauge. MemoryMax=6G remains the hard ceiling — this is strictly
# about removing the premature soft throttle.
MemoryMax=6G
MemoryHigh=infinity
SyslogIdentifier=studio-compositor

[Install]
WantedBy=default.target
