# systemd/units/studio-person-detector.service
[Unit]
Description=Studio Person Detector — YOLOv8n on per-camera snapshots
After=studio-compositor.service
Wants=studio-compositor.service
OnFailure=notify-failure@%n.service

[Service]
Type=simple
WorkingDirectory=/home/operator/projects/hapax-council
ExecStart=/home/operator/projects/hapax-council/.venv/bin/python -m agents.studio_person_detector
Restart=on-failure
RestartSec=10
Environment=PATH=/home/operator/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
MemoryMax=2G
SyslogIdentifier=studio-person-detector

[Install]
WantedBy=default.target
