# systemd/units/tmp-monitor.service
[Unit]
Description=Monitor /tmp usage and clean leaked files

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tmp-monitor.sh
