# systemd/units/video-retention.service
[Unit]
Description=Video recording retention — prune old segments

[Service]
Type=oneshot
ExecStart=/home/hapax/projects/hapax-council/scripts/video-retention.sh
Environment=HOME=/home/hapax
