# systemd/units/visual-layer-aggregator.service
[Unit]
Description=Visual Layer Aggregator — polls cockpit API and perception state for compositor overlay
After=cockpit-api.service hapax-voice.service
Wants=cockpit-api.service

[Service]
Type=simple
WorkingDirectory=/home/operator/projects/hapax-council
EnvironmentFile=-/run/user/1000/hapax-voice.env
ExecStart=/home/operator/projects/hapax-council/.venv/bin/python -m agents.visual_layer_aggregator
Restart=on-failure
RestartSec=10
Environment=PATH=/home/operator/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
MemoryMax=256M
SyslogIdentifier=visual-layer-aggregator

[Install]
WantedBy=default.target
