# systemd/units/visual-layer-aggregator.service
[Unit]
Description=Visual Layer Aggregator — polls logos API and perception state for compositor overlay
After=logos-api.service hapax-daimonion.service hapax-secrets.service
Requires=hapax-secrets.service
Wants=logos-api.service

[Service]
Type=simple
WorkingDirectory=/home/hapax/projects/hapax-council
EnvironmentFile=/run/user/1000/hapax-secrets.env
ExecStart=/home/hapax/projects/hapax-council/.venv/bin/python -m agents.visual_layer_aggregator
Restart=on-failure
RestartSec=10
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
MemoryMax=1G
SyslogIdentifier=visual-layer-aggregator

[Install]
WantedBy=default.target
