# systemd/units/weather-sync.timer
[Unit]
Description=Hourly weather sync

[Timer]
OnCalendar=*:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
