# systemd/units/weather-sync.service
[Unit]
Description=Weather sync — fetch current conditions for energy/mood correlation
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/hapax/.local/bin/uv run --directory /home/hapax/projects/hapax-council python -m agents.weather_sync
Environment=PATH=/home/hapax/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/hapax
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weather-sync
