# systemd/units/audio-recorder.service.d/archive-path.conf
[Service]
# LRR Phase 2 item 5 — FLAC output path migration.
# The base unit writes to its default path (``~/.cache/hapax-audio/``).
# This drop-in moves output under the stream-archive root so archival
# media live in one consistent hierarchy with HLS segments.
#
# Respected by audio-recorder.service when its entrypoint reads
# ``HAPAX_AUDIO_ARCHIVE_ROOT`` — see agents/audio_recorder/ for the
# env var resolution. The retention policy (R1-R5) treats these files
# identically to HLS segments once rotated.
Environment=HAPAX_AUDIO_ARCHIVE_ROOT=%h/hapax-state/stream-archive/audio
