# systemd/units/tabbyapi-hermes8b.service.d/gpu-pin.conf
[Service]
# GPU pin for the Hermes 3 8B second-instance TabbyAPI.
#
# **TODO (operator decision before deploy):** pick one of the two
# candidate allocations below and uncomment the chosen block. The
# existing tabbyapi.service (Qwen) drop-in at
# `systemd/units/tabbyapi.service.d/gpu-pin.conf` currently uses
# Option γ (`CUDA_VISIBLE_DEVICES=0,1`, both GPUs visible) which was
# authored for the now-deferred 70B layer split. For 8B parallel
# pivot, the Qwen drop-in should also be updated to pin to a single
# GPU — the two instances must NOT both see both GPUs.
#
# GPU inventory (under CUDA_DEVICE_ORDER=PCI_BUS_ID, empirically
# verified 2026-04-14T15:30Z by the real beta session — see
# ~/.cache/hapax/relay/context/2026-04-14-beta-phase-3-supplement-
# verified-preconditions.md § 1):
#
#   index 0 = RTX 5060 Ti (sm_120, 15.5 GiB total)
#     current residents: studio-compositor (~3-4 GB),
#                        hapax-imagination (~1-2 GB),
#                        hapax-daimonion STT + Kokoro
#   index 1 = RTX 3090    (sm_86,  23.6 GiB total)
#     current residents: tabbyapi (Qwen 3.5-9B, ~12 GB incl. KV cache),
#                        hapax-dmn (~3 GB)
#
# Hermes 3 8B EXL3 5.0bpw footprint: ~5 GB weights + ~4 GB KV cache
# at 8192 context = ~9 GB working set. Needs roughly that much free
# VRAM on whichever GPU it lands on.
#
# ---- Option A: Hermes on 3090 alongside Qwen + hapax-dmn ----
#
# Puts both TabbyAPI instances on the same GPU (3090). Total 3090
# footprint after deploy: Qwen ~12 GB + Hermes ~9 GB + hapax-dmn
# ~3 GB = ~24 GB. This is RIGHT AT the 23.6 GiB card limit — no
# headroom, first burst activation will OOM. Only viable if either
# Qwen KV cache or hapax-dmn is relocated first, or Hermes max_seq_len
# is dropped to 4096 (~2 GB KV cache) to buy headroom.
#
# Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
# Environment="CUDA_VISIBLE_DEVICES=1"
#
# ---- Option B: Hermes on 5060 Ti, Qwen stays on 3090 ----
#
# Pins Hermes to the 5060 Ti (15.5 GiB). Current 5060 Ti residents
# (compositor + imagination + daimonion audio stack) use roughly
# 5-7 GB, leaving ~8-10 GB free. Hermes 8B working set ~9 GB barely
# fits. Close to the edge but feasible if the operator verifies
# compositor + imagination idle VRAM before starting Hermes.
#
# This option REQUIRES updating the existing tabbyapi.service drop-in
# to pin Qwen to GPU 1 only (currently sees both GPUs). Otherwise
# Qwen may try to spill onto the 5060 Ti and clash with Hermes.
#
# Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
# Environment="CUDA_VISIBLE_DEVICES=0"
#
# ---- Operator choice ----
#
# Uncomment one of the two `Environment=` pairs above before running
# `scripts/install-units.sh`. Leaving both commented will cause the
# unit to start with the systemd default environment (no CUDA override)
# which is not a safe deployment state.
#
# Empirical verification before deploy (regardless of choice):
#   1. nvidia-smi baseline: confirm target GPU's free VRAM ≥10 GB
#   2. Load Hermes 8B, observe `nvidia-smi dmon -s um` for 60 s under
#      idle + under a synthetic load of 512-token generations at
#      max_seq_len 8192
#   3. If peak > 90% of target GPU VRAM, drop Hermes max_seq_len
#      to 4096 and re-test
#   4. Only then start using the *-hermes LiteLLM routes in anger
