# systemd/units/tabbyapi.service.d/gpu-pin.conf
[Service]
# LRR Phase 3 item 1 — partition reconciliation Option α → Option γ.
#
# Previous Option α (single-GPU, pre-Hermes 3):
#   Environment="CUDA_VISIBLE_DEVICES=1"
# Pinned TabbyAPI to GPU 1 (RTX 3090) only. That was the right
# placement for single-model Qwen3.5-9B (5.76 GiB) because the
# 3090's 24 GiB left 18 GiB of headroom for KV cache expansion.
#
# Option γ (dual-GPU, Hermes 3 ready):
#   Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
#   Environment="CUDA_VISIBLE_DEVICES=0,1"
# Makes BOTH GPUs visible to TabbyAPI so exllamav3 can layer-split
# Hermes 3 70B across them. Under Option γ the process sees:
#
#   cuda:0 → RTX 5060 Ti (sm_120, 15.5 GiB total, small overflow slice)
#   cuda:1 → RTX 3090    (sm_86,  23.6 GiB total, main weight slice)
#
# `CUDA_DEVICE_ORDER=PCI_BUS_ID` is REQUIRED before `CUDA_VISIBLE_DEVICES`
# on this rig. Without it, CUDA's default FASTEST_FIRST heuristic
# inverts the device index vs nvidia-smi and the 3090 would show as
# process index 0 instead of 1 — the Hermes 3 gpu_split ordering
# would then be wrong. Confirmed during W2.1 migration of
# studio-compositor. Beta re-verified at 2026-04-14T15:30Z via a
# direct PyTorch matmul test on both devices with this env, so the
# ordering `[5060 Ti @ index 0, 3090 @ index 1]` is empirically
# locked in (see ~/.cache/hapax/relay/context/2026-04-14-beta-
# phase-3-supplement-verified-preconditions.md § 1).
#
# The corresponding Hermes 3 tabbyAPI config.yml must declare
# `gpu_split: [2.75, 23.5]` (NOT `[23.5, 2.75]`) — 2.75 GiB on the
# 5060 Ti, 23.5 GiB on the 3090.
#
# Rollback to Option α: replace `CUDA_VISIBLE_DEVICES=0,1` with
# `CUDA_VISIBLE_DEVICES=1` and restart. The CUDA_DEVICE_ORDER line
# is safe to leave as-is under Option α.
Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
Environment="CUDA_VISIBLE_DEVICES=0,1"
