# systemd/units/tabbyapi-olmo.service.d/gpu-pin.conf
[Service]
# Queue #211 — OLMo parallel TabbyAPI GPU pinning.
#
# Goal: pin OLMo 3-7B to the RTX 3090 (cuda:1 under PCI_BUS_ID ordering)
# co-tenant with Qwen3.5-9B from the primary tabbyapi.service. The 3090
# has ~12 GB of headroom after Qwen's ~6 GB, which is more than enough
# for OLMo's ~5-6 GB at 5.0bpw + KV cache at 16384.
#
# CUDA_DEVICE_ORDER=PCI_BUS_ID + CUDA_VISIBLE_DEVICES=1 masks the 5060 Ti
# entirely so exllamav3 sees a single cuda:0 device that is actually the
# 3090. This matches the PCI_BUS_ID lockdown established for the primary
# tabbyapi.service gpu-pin.conf drop-in.
#
# Under this masking:
#   cuda:0 (process-visible) → RTX 3090 (host: cuda:1 under PCI_BUS_ID)
#   The 5060 Ti is invisible to this process.
#
# config.yml's `gpu_split: [0, 6]` is effectively single-GPU when only
# one device is visible — exllamav3 treats the ratio as a hint + loads
# on the single visible device.
#
# To move OLMo to the 5060 Ti instead (reason: freeing the 3090 for
# Hermes 70B reactivation — NOT currently planned per drop #62 §14):
#   Environment="CUDA_VISIBLE_DEVICES=0"  (5060 Ti under PCI_BUS_ID)
# But 5060 Ti has only ~10 GB free, tighter than the 3090.
Environment="CUDA_DEVICE_ORDER=PCI_BUS_ID"
Environment="CUDA_VISIBLE_DEVICES=1"
